# NX.Engine/Hive/Genomes/certbot/nxp/run/crontab
*   *   *   *   *  root /bin/bash /etc/nxp/run/renew-all > /dev/stdout
